# deploy/grimoire-bot.service
[Unit]
Description=Grimoire Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/grimoire
ExecStart=/opt/grimoire/grimoire-bot
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/grimoire
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
